# application/supervisord.conf
[program:browsertap]
command=/srv/browsertap/bin/cli -t slave